# birthday/mailsender.service
[Unit]
Description=Mail sendes

[Service]
Type=oneshot
ExecStart=/opt/sender/sendmail.py
